# recipes-graphics/wayland/weston-init/weston.service
[Unit]
Description=Weston Wayland Compositor
RequiresMountsFor=/run
After=synasdk-modules-load

[Service]
EnvironmentFile=-/etc/default/weston
StandardError=journal
PermissionsStartOnly=true
IgnoreSIGPIPE=no

ExecStart=/usr/bin/weston-start -v -e -- $OPTARGS
